# python-service/aiva-python.service
[Unit]
Description=AIVA Python Knowledge Service
After=network.target mysql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/etc/aiva-oai/python-service
Environment="PATH=/etc/aiva-oai/python-service/venv/bin"
ExecStart=/etc/aiva-oai/python-service/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 62002 --workers 4
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/aiva-python.log
StandardError=append:/var/log/aiva-python-error.log

[Install]
WantedBy=multi-user.target